# magi-archive-mcp-http.service
[Unit]
Description=Magi Archive MCP HTTP Server
After=network.target
Requires=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/magi-archive-mcp
Environment="RAILS_ENV=production"
EnvironmentFile=/home/ubuntu/magi-archive-mcp/.env.production
ExecStart=/home/ubuntu/.rbenv/shims/ruby /home/ubuntu/magi-archive-mcp/bin/mcp-server-http
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
